# deploy/ai-code-mother.service
[Unit]
Description=AI Code Mother Spring Boot Application
Documentation=https://github.com/vasc-language/ai-code-mother
After=network.target mysql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data

# 工作目录
WorkingDirectory=/var/www/ai-code-mother/backend

# 环境变量
Environment="JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64"
Environment="SPRING_PROFILES_ACTIVE=prod"
Environment="SERVER_PORT=8123"

# JVM 参数配置
Environment="JAVA_OPTS=-Xms512m -Xmx2048m -XX:+UseG1GC -XX:MaxGCPauseMillis=200 -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/var/log/ai-code-mother/heapdump.hprof"

# 启动命令
ExecStart=/usr/bin/java $JAVA_OPTS -jar /var/www/ai-code-mother/backend/ai-code-mother.jar --spring.profiles.active=prod

# 停止命令
ExecStop=/bin/kill -15 $MAINPID

# 重启策略
Restart=on-failure
RestartSec=10s

# 标准输出和错误输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-code-mother

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全配置
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
